# buffpy/scripts/buffbot.service
[Unit]
Description=CU Robotics Startup Service
ConditionPathExists=/home/cu-robotics/buff-code/buffpy

[Service]
Type=forking
Restart=always
User=cu-robotics
TimeoutStartSec=900
WorkingDirectory=/home/cu-robotics/buff-code
ExecStart=/bin/bash -c 'source /home/cu-robotics/buff-code/buffpy/buff.bash && run self.yaml'

[Install]
WantedBy=multi-user.target