# buffpy/scripts/buffbot.service
[Unit]
Description=CU Robotics Startup Service
ConditionPathExists=/home/cu-robotics/buff-code/buffpy

[Service]
User=cu-robotics
Type=forking
StandardOutput=console
WorkingDirectory=/home/cu-robotics/buff-code
ExecStart=/bin/bash -c 'source /home/cu-robotics/buff-code/buffpy/buff.bash && run self.yaml'

[Install]
#WantedBy=