# buffpy/scripts/buffbot.service
[Unit]
Description=CU Robotics Startup Service
ConditionPathExists=/home/cu-robotics/buff-code/buffpy

[Service]
Type=simple
RestartSec=5s
User=cu-robotics
WorkingDirectory=/home/cu-robotics/buff-code
ExecStart=/usr/bin/python3 /home/cu-robotics/buff-code/buffpy/bin/run self.yaml
Environment="PYTHONPATH=${PYTHONPATH}:/home/cu-robotics/.local/lib/python3.6/site-packages:/opt/ros/melodic/lib/python2.7/dist-packages"

[Install]
#WantedBy=